# systemd/db-backup.timer
[Unit]
Description=DB Backup Timer

[Timer]
Unit=db-backup.service
OnCalendar=monthly

[Install]
WantedBy=timers.target