# systemd/db-backup.timer
[Unit]
Description=EPG-Download Timer

[Timer]
OnCalendar=monthly

[Install]
WantedBy=multi-user.target