# systemd/dyndns-update.timer
[Unit]
Description=DynDns Update Timer

[Timer]
Unit=dyndns-update.service
OnCalendar=09:00

[Install]
WantedBy=timers.target