# systemd/epg-downloader.timer
[Unit]
Description=EPG-Download Timer

[Timer]
Unit=epg-downloader.service
OnCalendar=Mon,Wed,Sat 03:00

[Install]
WantedBy=timers.target