# systemd/mediathekview-update.timer
[Unit]
Description=Mediathek DB Update Timer

[Timer]
OnCalendar=04:00

[Install]
WantedBy=multi-user.target