# systemd/mediathekview-update.timer
[Unit]
Description=Mediathek DB Update Timer

[Timer]
Unit=mediathekview-update.service
OnCalendar=04:00

[Install]
WantedBy=timers.target